# src/etc/supervisord.d/10-system-timezone-wrapper.conf
[program:system-timezone-wrapper]
autorestart = false
autostart = true
command = /usr/sbin/system-timezone-wrapper --verbose
priority = 10
startretries = 0
startsecs = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
