# src/etc/supervisord.d/sshd-bootstrap.conf
[program:sshd-bootstrap]
priority = 5
command = /usr/sbin/sshd-bootstrap
autostart = %(ENV_SSH_AUTOSTART_SSHD_BOOTSTRAP)s
startsecs = 0
startretries = 0
autorestart = false
redirect_stderr = true
stdout_logfile = /var/log/secure
stdout_events_enabled = true