# etc/services-config/supervisor.d/sshd.conf
[program:sshd]
priority = 10
command = bash -c 'sleep 2 && /usr/sbin/sshd -D -e'
redirect_stderr = true
stdout_logfile = /var/log/secure
stdout_events_enabled = true