# etc/services-config/supervisor/supervisord.d/sshd-bootstrap.conf
[program:sshd-bootstrap]
priority = 5
command = bash -c 'touch /tmp/sshd-bootstrap.lock; env >> /etc/sshd-bootstrap.env; /usr/sbin/sshd-bootstrap && rm -f /tmp/sshd-bootstrap.lock'
startsecs = 0
startretries = 0
autorestart = false
redirect_stderr = true
stdout_logfile = /var/log/secure
stdout_events_enabled = true