# etc/services-config/supervisor/supervisord.d/sshd_bootstrap.conf
[program:sshd_bootstrap]
priority = 5
command = /etc/ssh-bootstrap
startsecs = 0
startretries = 0
autorestart = false
redirect_stderr = true
stdout_logfile = /var/log/secure
stdout_events_enabled = true